# overlay/files/etc/systemd/system/rescue.service.d/10-palette-debug.conf
# Same as the emergency drop-in, for the rescue.target path.
[Unit]
ConditionKernelCommandLine=palette.debug

[Service]
ExecStartPre=-/opt/spectrocloud/bin/collect-debug-bundle.sh --trigger=rescue
