# overlay/files/etc/systemd/system/kairos-agent.service.d/10-palette-debug.conf
# Best-effort: if the Kairos install agent runs as kairos-agent.service and
# fails, trigger the collector directly. Harmless if the unit does not exist
# under this name in a given framework build (systemd ignores the drop-in).
# The watchdog timer and emergency/rescue drop-ins are the framework-agnostic
# fallbacks that cover collection regardless of the install unit's name.
[Unit]
ConditionKernelCommandLine=palette.debug
OnFailure=palette-debug-collector.service
