# overlay/files/etc/systemd/system/emergency.service.d/10-palette-debug.conf
# When a debug install drops to emergency.target (e.g. blank screen, failed
# mount, install service crash), collect logs BEFORE the emergency shell.
# ExecStartPre with '-' is best-effort and never blocks the shell.
[Unit]
ConditionKernelCommandLine=palette.debug

[Service]
ExecStartPre=-/opt/spectrocloud/bin/collect-debug-bundle.sh --trigger=emergency
